# app/crontab
# Every minute
*/1 * * * * . /root/env.sh; /usr/local/bin/node /var/app/app/jobs/updateRanking.js >> /var/log/cron.log 2>&1
*/1 * * * * . /root/env.sh; /usr/local/bin/node /var/app/app/jobs/updateOverallRanking.js >> /var/log/cron.log 2>&1

# Every hour
0 * * * * . /root/env.sh; /usr/local/bin/node /var/app/app/jobs/incrementBalance.js >> /var/log/cron.log 2>&1

<<<<<<< Updated upstream
# Every Monday at 00:00
0 0 * * 1 . /root/env.sh; /usr/local/bin/node /var/app/app/jobs/closeCurrentRanking.js >> /var/log/cron.log 2>&1
=======
# At 00:00 on day-of-month 1
0 0 1 * * . /root/env.sh; /usr/local/bin/node /var/app/app/jobs/closeCurrentRanking.js --language=en >> /var/log/cron.log 2>&1
0 0 1 * * . /root/env.sh; /usr/local/bin/node /var/app/app/jobs/closeCurrentRanking.js --language=fr >> /var/log/cron.log 2>&1
>>>>>>> Stashed changes
